# mailhog.conf
[program:mailhog]
command=/usr/local/bin/MailHog -auth-file=/etc/auth.txt
autorestart=true
startretries=10